# manifest/etc/s6-overlay/s6-rc.d/roundcube/up
/etc/s6-overlay/script/setup_roundcube.sh
